# supervisord.conf
[program:ta]
directory=$HOME/me/ta
command=venv/bin/gunicorn main:application -c gunicorn.conf.py
user=$USER
autostart=true
autorestart=true
redirect_stderr=true